# backend/ops/db-retention/db-retention.cron
# Viru DB retention - run daily at 02:20 server time
# Install with: crontab backend/ops/db-retention/db-retention.cron
20 2 * * * /opt/viru-tracker/backend/ops/db-retention/run-db-retention.sh >> /opt/viru-tracker/backend/logs/db-retention-cron.log 2>&1
